# backup-via-usb-scripts/udev-rules/99-archiv-epg.rules
KERNEL=="sd?1", ATTRS{serial}=="574D43315431363934383331", RUN+="/usr/local/bin/mount_and_check_disk.sh $name /media/archive_epg"
