# backup-via-usb-scripts/udev-rules/99-backup-dd.rules
KERNEL=="sd?1", ATTRS{serial}=="5743415A4148353232393535", RUN+="/usr/local/bin/mount_and_check_disk.sh $name /media/backup-dd"
